# server/klip-uploader/deploy/klip-purge.timer
[Unit]
Description=Ejecuta la purga de capturas Klip cada hora

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
